# systemd/restic-retention@.timer
[Unit]
Description=Run restic retention (%i) weekly

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
