# systemd/restic-cache-clean.service
[Unit]
Description=Clean restic cache

[Service]
Type=oneshot
ExecStart=/usr/bin/restic cache --cleanup --max-age 30 --cache-dir /var/cache/restic
