# systemd/restic-retention.service
[Unit]
Description=Restic retention cleanup

[Service]
Type=oneshot
WorkingDirectory=/usr/local/bin/restic-ops
ExecStart=/usr/local/bin/restic-ops/bin/retention.sh
