# systemd/restic-backup.timer
[Unit]
Description=Run restic backup (%i) daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
