# systemd/restic-prune.timer
[Unit]
Description=Monthly Restic prune (heavy maintenance)

[Timer]
# Run on the first Sunday of the month (Days 1-7) at 04:00
OnCalendar=Sun *-*-01..07 04:00:00
Persistent=true
Unit=restic-prune.service

[Install]
WantedBy=timers.target
