# systemd/restic-cache-clean.timer
[Unit]
Description=Weekly restic cache cleanup

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
